# docker/supervisord-devel.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=python ../asilinks/manage.py shell_plus --notebook
